# arch-packages/client/clipboard-sync-client.install
#!/bin/bash

post_install() {
    echo "Enabling and starting clipboard-sync-client.service"
    systemctl enable clipboard-sync-client.service
    systemctl start clipboard-sync-client.service
}