# arch-packages/server/clipboard-sync-server.install
#!/bin/bash

post_install() {
    echo "Enabling and starting clipboard-sync-server.service"
    systemctl enable clipboard-sync-server.service
    systemctl start clipboard-sync-server.service
}